# deploy/ubuntu/twag-nimble@.service
[Unit]
Description=TWAG Nimble process (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/twag/twag.env
ExecStart=__APP_DIR__/deploy/ubuntu/run-nimble.sh
Restart=always
RestartSec=5
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
