# deploy/ubuntu/twag-sync-agent@.service
[Unit]
Description=TWAG sync agent (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/twag/twag.env
ExecStart=__APP_DIR__/deploy/ubuntu/run-sync-agent.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
Nice=10
CPUQuota=75%
IOSchedulingClass=best-effort
IOSchedulingPriority=7

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
